# hermes-meetily-watcher.service
[Unit]
Description=Hermes Meetily Export Processor
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/hermes-meetily-to-obsidian
ExecStart=/usr/bin/env python3 /root/hermes-meetily-to-obsidian/scripts/process_exporter.py --export-dir /root/meetily_exports --vault /root/Obsidian --cleanup-source keep
